# scripts/com.sina.brain-loop-autorun-v1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sina.brain-loop-autorun-v1</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@SG_ROOT@/scripts/brain_loop_launchd_entry_v1.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>@SG_ROOT@</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/opt/homebrew/bin:/usr/local/bin</string>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>USER</key>
    <string>@USER@</string>
    <key>SOURCEA_ROOT</key>
    <string>@HOME@/Projects/SourceA</string>
    <key>BRAIN_SG_ROOT</key>
    <string>@SG_ROOT@</string>
  </dict>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>@HOME@/.sina/logs/brain-loop-autorun-v1.out.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/.sina/logs/brain-loop-autorun-v1.err.log</string>
  <key>ThrottleInterval</key>
  <integer>300</integer>
</dict>
</plist>
